# agent/scripts/services/com.company.mdmagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.company.mdmagent</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/mdmagent</string>
        <string>-daemon</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/etc/mdmagent</string>
    
    <key>StandardOutPath</key>
    <string>/var/log/mdmagent/mdmagent-std.log</string>
    
    <key>StandardErrorPath</key>
    <string>/var/log/mdmagent/mdmagent-error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>UserName</key>
    <string>root</string>
    
    <key>GroupName</key>
    <string>wheel</string>
    
    <key>InitGroups</key>
    <true/>
    
    <key>LaunchOnlyOnce</key>
    <false/>
    
    <key>AbandonProcessGroup</key>
    <true/>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>ExitTimeOut</key>
    <integer>30</integer>
</dict>
</plist>